# crontab/crontab
*/10 * * * * /root/script-docker-java-trading.sh